# deploy/pushgo-gateway.service
[Unit]
Description=PushGo Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=pushgo
Group=pushgo
WorkingDirectory=/var/lib/pushgo-gateway

ExecStart=/usr/local/bin/pushgo-gateway
Restart=on-failure
RestartSec=3

Environment="PUSHGO_HTTP_ADDR=127.0.0.1:6666"
Environment="PUSHGO_TOKEN="
Environment="PUSHGO_GATEWAY_URL=https://gateway.pushgo.dev"
Environment="MAX_CONCURRENT=200"
Environment="DATA_PATH=/var/lib/pushgo-gateway/data"
Environment="PUSHGO_DB_URL="

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/pushgo-gateway

[Install]
WantedBy=multi-user.target
